# .config/systemd/user/urxvtd.service
[Unit]
Description=URxvt Terminal Daemon

[Service]
ExecStart=/usr/bin/urxvtd -q -o

[Install]
WantedBy=default.target
